# conf/supervisor.ini
[program:app]
command = uwsgi /work/conf/uwsgi.ini
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /data/log/supervisor.log
stopasgroup = true
killasgroup = true
